# server-ce/runit/filestore-sharelatex/run
#!/bin/bash
exec /usr/bin/node /var/www/sharelatex/filestore/app.js >> /var/log/sharelatex/filestore.log 2>&1
